# meta-arago-distro/recipes-core/arago-gpl-notice/arago-gpl-notice/gplv3-notice.service
[Unit]
Description=Print notice about GPLv3 packages

[Service]
Type=oneshot
ExecStart=/usr/bin/gplv3-notice
WorkingDirectory=/usr/bin
StandardOutput=tty

[Install]
WantedBy=multi-user.target
